# deploy/unly.service
[Unit]
Description=Unly Personal AI Agent Platform
Documentation=https://github.com/oblivisheee/unly
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=unly
Group=unly
WorkingDirectory=__UNLY_WORKING_DIRECTORY__
ExecStart=__UNLY_EXECUTABLE__ start --config __UNLY_CONFIG_PATH__
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10

# Environment file for secrets (not version-controlled)
EnvironmentFile=-/etc/unly/unly.env

# Logging — integrates with journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unly

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=200%

[Install]
WantedBy=multi-user.target
