# deploy/unly.service
[Unit]
Description=Unly Personal AI Agent Platform
Documentation=https://github.com/oblivisheee/unly
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=unly
Group=unly
WorkingDirectory=/opt/unly
ExecStart=/opt/unly/bin/unly start --config /opt/unly/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10

# Environment file for secrets (not version-controlled)
EnvironmentFile=-/etc/unly/unly.env

# Logging — integrates with journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unly

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# ProtectHome must remain false when the workspace lives inside the user home
# directory (the default: ~/.unly).  Enabling read-only breaks database writes.
ProtectHome=false
ReadWritePaths=/opt/unly
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=200%

[Install]
WantedBy=multi-user.target
